# infra/compute-engine/robocop-control-agent.service
[Unit]
Description=Robocop Control Agent
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/opt/robocop/src/services/control_agent
Environment=NODE_ENV=production
Environment=PORT=8080
Environment=DEVICE_ID=rover-dev-001
Environment=COMMAND_TIMEOUT_MS=500
ExecStart=/usr/bin/node dist/main.js
Restart=always
RestartSec=5
User=robocop
Group=robocop

[Install]
WantedBy=multi-user.target
